# deploy/supervisord.ini
[program:craigle_worker]
directory=/srv/craigle
command=python3 manage.py rqworker default
numprocs=25
process_name=%(program_name)s_%(process_num)s
autostart=true
autorestart=true
user=apache
environment=PATH="/srv/craigle/venv/bin"
